# systemd/user/mpd.service.d/override.conf
[Service]
ExecStartPre=/bin/mkdir -p %t/mpd %C/mpd
Restart=on-abnormal