# sv/mpd/run
#!/bin/sh
install -d -m0750 ~/.run/mpd
exec mpd --no-daemon